# deployment/systemd/runtime-ai-agent.service
[Unit]
Description=Runtime AI Security Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/runtime-ai-agent \
  --cloud-endpoint=https://api.runtime-ai-security.example.com \
  --mode=permissive \
  --policy-path=/etc/runtime-ai-security/policy.json

# Restart configuration
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=runtime-ai-agent

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/runtime-ai-security

[Install]
WantedBy=multi-user.target
